# encfs
#!/sbin/runscript

depend()
{
	need localmount
	use netmount nfsmount
}

start()
{
	ebegin "Mount users encrypted dirs"
	eend 0
}

stop()
{
	ebegin "Deactivating encrypted dirs"
	eend 0
}
